# runit/powertune/run
#!/bin/sh -e

powertop --auto-tune

exec chpst -b powertune pause
